# 30-Resources/Hermes/Journal-KPI/journal-kpi-watcher.service
[Unit]
Description=Obsidian Journal KPI Watcher
After=network.target
# Restart if it crashes
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=simple
EnvironmentFile=/home/mike/.hermes/.env
ExecStart=/home/mike/.hermes/scripts/journal_watcher.sh
Restart=on-failure
RestartSec=10
StandardOutput=append:/home/mike/.hermes/logs/journal_watcher.log
StandardError=append:/home/mike/.hermes/logs/journal_watcher.log

[Install]
WantedBy=default.target
